# deploy/compute-wrapped.cron.example
# Cron job: rebuild Plex Wrapped cache
#
# Install:
#   sudo cp deploy/compute-wrapped.cron.example /etc/cron.d/plex-wrapped-compute
#
# Runs daily at 03:15. Adjust schedule as needed.
# Uses the production database (not USE_TEST_DATABASE).
# Year comes from WRAPPED_YEAR in .env (no --year flag needed).

SHELL=/bin/bash
PATH=/usr/local/bin:/usr/bin:/bin

15 3 * * * plexwrapped cd /opt/plex-wrapped && .venv/bin/python scripts/compute_wrapped.py >> /var/log/plex-wrapped-compute.log 2>&1
